# com.horos.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.horos.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Volumes/PACS/Backup/horos_backup_export.py</string>
  </array>

  <!-- Executa a cada 10 minutos -->
  <key>StartInterval</key>
  <integer>600</integer>

  <key>RunAtLoad</key>
  <true/>

  <!-- Logs do launchd (secundários). O script já grava logs rotacionados em /Volumes/PACS/Backup/logs -->
  <key>StandardOutPath</key>
  <string>/tmp/horos_backup_export.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/horos_backup_export.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>LANG</key><string>en_US.UTF-8</string>
    <key>LC_ALL</key><string>en_US.UTF-8</string>
  </dict>
</dict>
</plist>